# app/home/config/systemd/user/auto-update-gnome-theme.path
[Unit]
Description=watch file

[Path]
PathModified=/var/lib/gnome-night-shift/is-day-or-night

[Install]
WantedBy=default.target

