# app/home/config/systemd/user/auto-update-gnome-theme.path
[Unit]
Description=watch file

[Path]
PathModified=%h/.local/state/gnome-mode-shift/is-day-or-night

[Install]
WantedBy=default.target

