# app/home/config/systemd/user/gnome-mode-shift.service
[Unit]
Description=Sets gsettings prefer-<light|dark>
Wants=network-online.target display-manager.service
After=network-online.target graphical-session.target

[Service]
StateDirectory=%p
Type=oneshot
WorkingDirectory=/usr/lib/gnome-mode-shift
ExecStart=/usr/lib/gnome-mode-shift/bin/get-sunrise-sunset.sh
