# app/home/config/systemd/user/stop-mode-shift-timer.service
[Unit]
Description=Stops gnome-mode-shift timer on failure

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user stop gnome-mode-shift.timer
